# deploy/wisp-monitor.service
# WISP Monitor — edge probe (central-brain mode).
#
# The edge no longer runs a local FSM/alerting/dashboard: it only probes and reports
# raw results to central, which does all the detecting and paging. Gives auto-start
# on boot and auto-restart on crash. Adjust WorkingDirectory/User to your install.
#
#   sudo cp deploy/wisp-monitor.service /etc/systemd/system/
#   sudo systemctl daemon-reload && sudo systemctl enable --now wisp-monitor
#
[Unit]
Description=WISP Monitor (edge probe, central-brain mode)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Where the repo lives on the on-prem box:
WorkingDirectory=/opt/wisp
# Use the venv python so icmplib/httpx are importable:
ExecStart=/opt/wisp/.venv/bin/python /opt/wisp/apps/daemon/main.py
Restart=always
RestartSec=3
Environment=PYTHONUNBUFFERED=1
# REQUIRED — point this box at your central server and identify it. Ask your
# central operator for the URL + token; tenant/node id are yours to choose
# (node_id defaults to the hostname if left unset).
Environment=WISP_CENTRAL_BRAIN=1
Environment=WISP_CENTRAL_URL=https://central.example.net
Environment=WISP_CENTRAL_TOKEN=changeme
Environment=WISP_TENANT_ID=changeme
# Environment=WISP_NODE_ID=this-box
# Detection speed: poll cadence × down_consecutive (3) ~= time-to-DOWN, enforced by
# central's FSM. Central also sends a fast-confirm "recheck" hint for a device that
# looks suspect, so wall-clock detection is usually much faster than this alone.
Environment=WISP_POLL_INTERVAL_S=20
# Other tunables (thresholds, retry cadence) can go here as env too — see config.py.
# The ICMP prober uses unprivileged ping sockets — no cap_net_raw needed; just enable
# the kernel ping group once: sudo sysctl -w net.ipv4.ping_group_range="0 2147483647"
# (persist in /etc/sysctl.d/). Then it can run as a non-root user:
# User=wisp
# Group=wisp

[Install]
WantedBy=multi-user.target
